# shadowsocks/sslocal.conf
[program:sslocal]
command=/bin/bash -c 'sslocal -c /etc/sslocal.json'
user=root
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/tmp/sslocal.log
stderr_logfile=/tmp/sslocal.err